# deploy/yeti-updater.service
[Unit]
Description=Yeti auto-updater
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/yeti/deploy/deploy.sh
