# deploy/yeti-updater-trigger.path
[Unit]
Description=Yeti manual update-check trigger

[Path]
PathExists=/home/yeti/.yeti/update-check-requested
Unit=yeti-updater.service

[Install]
WantedBy=paths.target
