# deploy/yeti.service
[Unit]
Description=Yeti - GitHub Automation Service
After=network.target

[Service]
Type=simple
User=yeti
Group=yeti
WorkingDirectory=/opt/yeti
EnvironmentFile=-/home/yeti/.yeti/env
ExecStart=/usr/bin/node /opt/yeti/dist/main.js
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=330

[Install]
WantedBy=multi-user.target
